# Logstash/supervisord.conf
[supervisord]
nodaemon=true

[program:logstash]
command=/tmp/logstash.sh
